# etc-systemd-system/-.automount
[Unit]
Description=/
Before=local-fs.target systemd-journald.socket
After=system.slice
Wants=system.slice

[Automount]
Where=/
