# etc-systemd-system/auto-chmod-uga\x2br-torrents@.path
[Unit]
Description=Make Readable All Torrent Files Directory Watch
After=transmission.service
ConditionPathExists=/%I

[Path]
PathChanged=/%I

[Install]
WantedBy=transmission.service
