# etc-systemd-system/media-ed\x2dusers.mount
[Unit]
Description=/media/ed-users
Before=remote-fs.target
After=NetworkManager-wait-online.service
Wants=NetworkManager-wait-online.service

[Mount]
UMask=0000
What=//edslaptop/Users
Where=/media/ed-users
Type=cifs
Options=rw,iocharset=utf8,sec=ntlm,credentials=/home/.smbcredentials
TimeoutSec=15

[Install]
WantedBy=remote-fs.target
